# supervisor/supervisord-bitcoind.conf
[program:bitcoind]
command=start-bitcoind
numprocs=1
autostart=true
autorestart=true
