# supervisor/supervisord-bitcore.conf
[program:bitcore]
nodaemon=true
command=start-bitcore
numprocs=1
autostart=true
autorestart=true
